# src/systemd/rrh-daily-check.service
[Unit]
Description="RRH daily outdate check"

[Service]
Type=oneshot
ExecStart=/usr/bin/rrh-auto-check